# deploy/nas/agent-executors/com.zz-agent.mimo-executor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>ZZ_BASE_URL</key>
		<string>http://192.168.31.119:18080/agent</string>
		<key>ZZ_EXECUTOR_INTERVAL</key>
		<string>30</string>
		<key>MIMO_AGENT_KEY</key>
		<string>zzk_REVOKED_PLACEHOLDER</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.zz-agent.mimo-executor</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Developer/CommandLineTools/usr/bin/python3</string>
		<string>/Users/z/.zz-agent/mimo-worker-executor-wrapper.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/z/.zz-agent/mimo-executor.launchd.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/z/.zz-agent/mimo-executor.launchd.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/z</string>
</dict>
</plist>
